# web/supervisord.conf
; supervisord config for the frontend sidecar.
;
; nginx is the primary process — Container Apps health checks hit
; http://127.0.0.1:8081/ which nginx serves. The cgroup reporter is
; advisory: if it crashes the SPA still works, so we mark it
; autorestart=true but exitcodes=0,1,2 keep the supervisor from cycling
; the entire pod over telemetry.
;
; Both children write logs to stdout/stderr so they end up in Container
; Apps log streams.

[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10

[program:cgroup-reporter]
command=/usr/bin/python3 /usr/local/bin/elb-cgroup-reporter %(ENV_SIDECAR_NAME)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
exitcodes=0,1,2
startretries=3
priority=20
environment=LOG_LEVEL="INFO"
